# tools/lib32_ubuntu/etc/profile.d/lib32_paths.sh
# add 32-bit library to search path
export LD_LIBRARY_PATH=/lib32:"$LD_LIBRARY_PATH"
export GTK_PATH=/lib32/gtk-2.0:"$GTK_PATH"
export LIBGL_DRIVERS_PATH=/lib32/dri:"LIBGL_DRIVERS_PATH"
